# roles/jenkins/templates/jenkins.service.j2
[Unit]
Description=Jenkins

[Service]
ExecStart=/usr/bin/java -jar /opt/jenkins/jenkins.war --httpPort=8081
Restart=always
User=ubuntu

[Install]
WantedBy=multi-user.target
